# ufw-notifier.service
[Unit]
Description=UFW Real-Time Notifier Daemon
After=network.target ufw.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ufw-notifier
ExecStart=/usr/bin/python3 /opt/ufw-notifier/ufw_watcher.py
Restart=always
RestartSec=5

# Hardening / Sandboxing
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
