# ufw-notifier.service
[Unit]
Description=UFW Real-Time Notifier Daemon
After=network.target ufw.service

[Service]
Type=simple
User=root
# We run as root or a user with permissions to read /var/log/ufw.log
# In a real environment, you'd create a dedicated user and grant read access.
WorkingDirectory=/opt/ufw-notifier
ExecStart=/usr/bin/python3 /opt/ufw-notifier/ufw_watcher.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
